# templates/pihole-update-gravity.service.j2
[Unit]
Description=Execute the pihole update gravity process

[Service]
ExecStart=/usr/bin/docker container exec {{ container_name }} pihole updateGravity

[Install]
WantedBy=default.target
